# .config/systemd/user/atuin-daemon.service
[Unit]
Description=atuin shell history daemon

[Service]
ExecStart=atuin daemon
Type=simple
Restart=always

[Install]
WantedBy=default.target
