# .config/systemd/user/background-daemons.slice
[Unit]
Description=User Slice for daemons
Documentation=man:systemd.special(7)
