# .config/systemd/user/background-containers.slice
[Unit]
Description=User Slice for Containers
Documentation=man:systemd.special(7)
