# .config/systemd/user/desktop.slice
[Unit]
Description=User Slice for Desktop Environment
Documentation=man:systemd.special(7)
