# .config/systemd/user/desktop-bar.target
[Unit]
Description=Target to signal that a bar is available
Wants=waybar.service
