# .config/systemd/user/restore-workspaces.service
[Unit]
Description=Restore saved workspaces
After=xrandr.service

[Service]
ExecStart=%h/.i3/scripts/restore-workspaces.sh
ExecStartPost=autotiling & diswon

[Install]
WantedBy=default.target
