# .config/systemd/user/run-mount-scripts.service
[Unit]
Description=Run mount.sh scripts found in /run/media/

[Service]
Type=exec
Environment="DISPLAY=:0" "PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin" "HOME=%h"
ExecStartPost=%h/.local/bin/umount_scripts
ExecStart=%h/.local/bin/mount_scripts
StandardOutput=journal

[Install]
WantedBy=default.target
