# .config/systemd/user/wallpaper.service
[Unit]
Description=Choose a random wallpaper and update betterlockscreen cache

[Service]
Type=oneshot
ExecStart=/usr/bin/betterlockscreen -u %h/.wallpaper --fx dimblur
ExecPostStart=sh -c "killall -9 hyprpaper && hyprpaper"
StandardOutput=journal

[Install]
WantedBy=default.target
