# .config/systemd/user/waybar.service
[Unit]
Description=Polybar

[Service]
Type=simple
StandardOutput=journal
ExecStart=%h/.config/waybar/launch.sh

[Install]
WantedBy=default.target
