# .config/systemd/user/x-session.target
[Unit]
Description=X11 session
Documentation=man:systemd.special
BindsTo=user-graphical-session.target
AssertEnvironment=DISPLAY
