# config/systemd/user/xrandr.service
[Unit]
Description=Configure displays

[Service]
ExecStartPre=/bin/sleep 2
Type=oneshot
ExecStart=/home/borhane/.files/layouts/main-setup.sh
ExecStartPost=/usr/bin/bash /home/borhane/.files/i3/scripts/bg.sh

[Install]
WantedBy=default.target
