# .config/systemd/user/gvfs-daemon.service.d/override.conf
[Service]
Slice=background-daemons.slice
